# com.user.loginscript.plist
<plist version="1.0">
 <dict>
   <key>Label</key>
   <string>customlogin-daemon</string>
   <key>RunAtLoad</key>
   <true />
   <key>Program</key>
   <string>/Users/emeshuris/git/dotsettings/custom_login</string>
  </dict>
</plist>
